# supervisor.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn NoInventoryProject.wsgi  --bind 0.0.0.0:80
directory=~/NoInventory/manage.py
user=root
autostart=true
autorestart=true
redirect_stderr=true
