# supervisor.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn practica4.wsgi  --bind 0.0.0.0:8000
directory=~/DAIproducion/manage.py
user=root
autostart=true
autorestart=true
redirect_stderr=true
